# package/base-files/files/etc/uci-defaults/100-qiaojie-wifi-5g
# 1. 创建网络接口
uci set network.wwan=interface
uci set network.wwan.proto='dhcp'

# 2. 配置无线参数
uci set wireless.wifinet2=wifi-iface
uci set wireless.wifinet2.device='radio1'
uci set wireless.wifinet2.mode='sta'
uci set wireless.wifinet2.network='wwan'
uci set wireless.wifinet2.ssid='WIFI_5G'
uci set wireless.wifinet2.encryption='sae'
uci set wireless.wifinet2.key='88886666'
# 针对某些驱动的兼容性设置
uci set wireless.radio1.cell_density='0'
uci set wireless.wifinet2.ocv='0'
# 将 wwan 接口加入到包含 wan 的防火墙区域中
# 假设第一个 zone 是 lan (index 0), 第二个是 wan (index 1)
uci add_list firewall.@zone[1].network='wwan'

# 如果你不确定索引，可以使用更安全的方法（找到 name 为 wan 的 zone）
# 这里我们直接增加允许 WAN 访问本地的规则
uci add firewall rule
uci set firewall.@rule[-1].name='ALLOW_LOCAL_ACCESS'
uci set firewall.@rule[-1].src='wan'
uci set firewall.@rule[-1].target='ACCEPT'