# package/base-files/files/etc/uci-defaults/100-qiaojie-wifi-5g
uci del firewall.cfg02dc81.network
uci add_list firewall.cfg02dc81.network='lan'
uci del firewall.cfg03dc81.network
uci add_list firewall.cfg03dc81.network='wan'
uci add_list firewall.cfg03dc81.network='wan6'
uci add_list firewall.cfg03dc81.network='wwan
'
uci add_list firewall.cfg03dc81.network='WIFI_CLIENT'
uci del firewall.docker.network
uci add_list firewall.docker.network='docker'
# /etc/config/network
uci set network.WIFI_CLIENT=interface
uci set network.WIFI_CLIENT.proto='dhcp'
# /etc/config/wireless
uci set wireless.wifinet2=wifi-iface
uci set wireless.wifinet2.device='radio1'
uci set wireless.wifinet2.mode='sta'
uci set wireless.wifinet2.network='WIFI_CLIENT'
uci set wireless.wifinet2.ssid='WIFI_5G'
uci set wireless.wifinet2.encryption='sae'
uci set wireless.wifinet2.key='88886666'
uci set wireless.radio1.cell_density='0'
uci set wireless.wifinet2.ocv='0'


uci add firewall rule # =cfg0f92bd
uci set firewall.@rule[-1].src='wan'
uci set firewall.@rule[-1].name='ALLOW_ACCESS_LOCAL'
uci set firewall.@rule[-1].target='ACCEPT'